# data/calls-daemon.service.in
[Unit]
Description=A phone dialer and call handler (daemon mode)
PartOf=gnome-session.target
Requisite=gnome-session-initialized.target
After=gnome-session-initialized.target

[Service]
BusName=org.gnome.Calls
ExecStart=@bindir@/gnome-calls --daemon
Type=dbus

[Install]
WantedBy=gnome-session.target
